# deploy/systemd/experience-api.service
[Unit]
Description=Experience Accumulation API (FastAPI)
After=network.target

[Service]
Type=simple
# 线上以 root 运行；如需低权用户部署，改为部署用户并调整目录属主
User=root
Group=root
WorkingDirectory=/opt/experience/backend
EnvironmentFile=/opt/experience/backend/.env
ExecStart=/opt/experience/backend/.venv/bin/uvicorn app.main:app \
    --host 127.0.0.1 \
    --port 8094 \
    --workers 2 \
    --proxy-headers \
    --forwarded-allow-ips="127.0.0.1"

Restart=always
RestartSec=3

# 安全加固
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/experience/backend/data /opt/experience/backend/storage /opt/experience/backend/backups

[Install]
WantedBy=multi-user.target
