# deploy/systemd/experience-api.service
[Unit]
Description=Experience Accumulation API (FastAPI)
After=network.target

[Service]
Type=simple
# 修改为实际部署用户与路径
User=www-data
Group=www-data
WorkingDirectory=/var/www/experience/backend
EnvironmentFile=/var/www/experience/backend/.env
ExecStart=/var/www/experience/backend/.venv/bin/uvicorn app.main:app \
    --host 127.0.0.1 \
    --port 8000 \
    --workers 2 \
    --proxy-headers \
    --forwarded-allow-ips="127.0.0.1"

Restart=always
RestartSec=3

# 安全加固
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/experience/backend/data /var/www/experience/backend/storage /var/www/experience/backend/backups

[Install]
WantedBy=multi-user.target
